# ahu_dashboard_pizero/rpi_ota_updater/ahu-ota-updater.service
[Unit]
Description=AHU Dashboard OTA Updater Service
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=almed
Group=almed
WorkingDirectory=/home/almed/Documents/almed_ahu/ahu_dashboard/rpi_ota_updater
ExecStart=/usr/bin/python3 /home/almed/Documents/almed_ahu/ahu_dashboard/rpi_ota_updater/rpi_ota_updater.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables (can be overridden in /etc/default/ahu-ota-updater)
Environment="MQTT_BROKER=localhost"
Environment="MQTT_PORT=1883"
Environment="MQTT_USERNAME=ahu_user"
Environment="MQTT_PASSWORD=ahu_pass_2024"
Environment="DASHBOARD_DIR=/home/almed/Documents/almed_ahu"
Environment="FLUTTER_PI_SERVICE=ahu-dashboard"
Environment="GITHUB_REPO_OWNER=ESPUpdaterzaid"
Environment="GITHUB_REPO_NAME=almed-rpi-dashboard"

# Optional: Load additional environment from file
EnvironmentFile=-/etc/default/ahu-ota-updater

[Install]
WantedBy=multi-user.target

